# scripts/titan-dispatcher.service
[Unit]
Description=Titan Dispatcher (safe starter wrapper)
After=network.target redis.service

[Service]
Type=simple
User=titan
WorkingDirectory=/home/titan/TitanNetwork
ExecStart=/bin/bash -c /home/titan/TitanNetwork/scripts/start_dispatcher_safe.sh
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

# Harden service
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
KillMode=control-group

[Install]
WantedBy=multi-user.target
